# deploy/jayendra-online.service
[Unit]
Description=Next.js for jayendra.online
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/jayendra.online

Environment=NODE_ENV=production
Environment=PORT=3001
EnvironmentFile=-/root/jayendra.online/.env.local
Environment=PATH=/usr/bin:/bin

ExecStart=/usr/bin/npm run start

Restart=always
RestartSec=5
TimeoutStopSec=30
LimitNOFILE=65536
MemoryMax=1200M
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jayendra-online

[Install]
WantedBy=multi-user.target
